# comenv/i3/etc/supervisor/conf.d/nx.conf
[program:nx]
command=/opt/nxserver.sh
user=root
autorestart=true
priority=400
